# modules/mesaides/templates/openfisca.conf.erb
start on runlevel [2345]
stop on runlevel [!2345]
respawn

setuid main
setgid main

chdir /home/main/mes-aides-ui/openfisca

exec <%= @venv_dir %>/bin/gunicorn api --config config.py
